# .config/systemd/user/backup.target
[Unit]
Description=Backup Timer Target
StopWhenUnneeded=yes

[Install]
WantedBy=timers.target
