# .config/systemd/user/backup.timer
[Unit]
Description=Timer for backup jobs

[Timer]
OnCalendar=daily
Unit=backup.target
