# .config/systemd/user/borg-nas.timer
[Unit]
Description=Run borg

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
