# .config/systemd/user/borg-cloud-personal.timer
[Unit]
Description=Run borg

[Timer]
OnCalendar=Mon..Fri *-*-* 10:30:00

[Install]
WantedBy=timers.target
