# .config/systemd/user/sync-to-nas.timer
[Unit]
Description=Run sync-to-nas

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
