# .config/systemd/user/sync-to-nas.timer
[Unit]
Description=Run sync-to-nas

[Timer]
OnCalendar=Mon..Fri *-*-* 10:00:00

[Install]
WantedBy=timers.target
